# packaging/systemd/forge-bridge.service
[Unit]
Description=forge-bridge MCP server + Artist Console (port 9996 + MCP HTTP port 9997)
Documentation=https://github.com/cnoellert/forge-bridge/blob/main/docs/INSTALL.md
Requires=forge-bridge-server.service
After=forge-bridge-server.service network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=__SUDO_USER__
Group=__SUDO_USER__
EnvironmentFile=/etc/forge-bridge/forge-bridge.env
ExecStart=/home/__SUDO_USER__/miniconda3/envs/forge/bin/python -m forge_bridge --transport streamable-http --mcp-port 9997
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
WorkingDirectory=/home/__SUDO_USER__

[Install]
WantedBy=multi-user.target
